# TP3/yt.service
[Unit]
Description=My youtube downloader

[Service]
ExecStart=bash /srv/yt/yt-v2.sh

[Install]
WantedBy=multi-user.target
